# trimble-agentic-docs-mcp/examples/trimble-agentic-docs-refresh.timer.example
[Unit]
Description=Run Trimble Agentic docs refresh weekly

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=30m
Unit=trimble-agentic-docs-refresh.service

[Install]
WantedBy=timers.target
